# 99-input.rules
KERNEL=="input[0-9]*", RUN+="/usr/lib/udev/bluetooth"
