# radio.service
[Unit]
Description=Radio

[Service]
Type=simple
Restart=always
PrivateTmp=true
ExecStart=/bin/bash /home/alarm/radioStart.sh
